# services/surg-db-frontend.service
[Unit]
Description=Surgical Database Frontend
After=network.target

[Service]
Type=simple
User=root
WorkingDirectory=/root/frontend
Environment="PATH=/usr/local/bin:/usr/bin:/bin"
ExecStart=/usr/bin/npm run dev -- --host 0.0.0.0 --port 3000
Restart=always
RestartSec=10
StandardOutput=append:/root/.tmp/frontend.log
StandardError=append:/root/.tmp/frontend.log

[Install]
WantedBy=multi-user.target
